# update-motd.d/99-footer
#!/bin/sh
exec figlet -f standard $(hostname)
printf "\n"

